# ha-kiosk-google-calendar-sync/etc/systemd/system/ha-kiosk-google-calendar-sync.service
[Unit]
Description=Calendar sync daemon for the kiosk app
After=network-online.target
Wants=network-online.target

[Service]
User=kiosk
Group=kiosk
ExecStart=/opt/kiosk/bin/ha-kiosk-google-calendar-sync --config /etc/kiosk/daemon-config.json
# Shared with ha-kiosk-weather-sync.service (same RuntimeDirectory=kiosk —
# ha-kiosk the app reads both daemons' snapshots from one well-known dir).
# Preserve=yes so a restart of *this* unit doesn't blow away the other
# daemon's snapshot/lock files sitting in the same directory; still starts
# fresh at actual boot, since /run itself is a new tmpfs then.
RuntimeDirectory=kiosk
RuntimeDirectoryPreserve=yes
RuntimeDirectoryMode=0750
LogsDirectory=ha-kiosk-google-calendar-sync
LogsDirectoryMode=0750
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
